# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true

[program:php-fpm7.4]
command=php-fpm7.4 -y /etc/php/7.4/fpm/pool.d/www.conf -F
