# supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g 'daemon off;'
autorestart=true

[program:php5-fpm]
command=php5-fpm -y /etc/php5/fpm/pool.d/www.conf -F